# i3batwarn.timer
[Unit]
Description=i3 Battery Warning Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
